# deploy/device/etc/systemd/system/physicar-myapp.service
[Unit]
Description=PhysiCar MyApp Service
After=network-online.target
Wants=network-online.target
# Only run (and only create run.log) when the student has actually deployed an
# app. Without this the service starts on every boot even with no run.sh and
# leaves an empty run.log behind. The web UI writes myapp/run.sh then restarts
# the unit (routers/myapp.py), so the condition is re-evaluated at that point.
ConditionPathExists=/opt/physicar/userdata/myapp.sh

[Service]
Type=simple
User=physicar
Group=physicar
WorkingDirectory=/opt/physicar
Environment=PORT=5000 PYTHONUNBUFFERED=1
# truncate: lets systemd (root) create+truncate the log itself, so no ExecStartPre
# is needed and there is no physicar-vs-root ownership clash on first start.
# (The myapp/ dir itself is guaranteed by routers/myapp.py before run.sh exists.)
StandardOutput=truncate:/opt/physicar/userdata/myapp.log
StandardError=append:/opt/physicar/userdata/myapp.log
ExecStart=/bin/bash -ic 'exec bash /opt/physicar/userdata/myapp.sh'
Restart=on-failure
RestartSec=2
KillMode=control-group

[Install]
WantedBy=multi-user.target
