# deploy/sim/supervisord.conf
[supervisord]
nodaemon=false
; Keep Python bytecode caches out of the workspace — no __pycache__ in student
; folders, and the watchdog reloader stops reacting to cache-file creation
environment=PYTHONPYCACHEPREFIX="/tmp/pycache"
logfile=/tmp/supervisord.log
logfile_maxbytes=5MB
logfile_backups=2
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:xvfb]
command=Xvfb :1 -screen 0 800x600x24 -ac +extension GLX +render -noreset
autorestart=true
priority=200
stdout_logfile=/dev/null
stderr_logfile=/dev/null
redirect_stderr=true

[program:openbox]
command=openbox
environment=DISPLAY=":1"
autorestart=true
priority=210
depends_on=xvfb
stdout_logfile=/dev/null
stderr_logfile=/dev/null
redirect_stderr=true

[program:tint2]
command=tint2
environment=DISPLAY=":1"
autorestart=true
priority=215
depends_on=openbox
stdout_logfile=/dev/null
stderr_logfile=/dev/null
redirect_stderr=true

[program:xterm]
command=/bin/bash -c "xhost +SI:localuser:physicar &>/dev/null || true; exec xterm -geometry 80x24+80+50 -fa Monospace -fs 10 -bg black -fg white -title Terminal -e 'cd ~ && bash -l'"
environment=DISPLAY=":1",HOME="/home/physicar"
autorestart=false
startsecs=0
priority=217
depends_on=openbox
stdout_logfile=/dev/null
stderr_logfile=/dev/null
redirect_stderr=true

[program:x11vnc]
command=x11vnc -display :1 -forever -nopw -shared -rfbport 5901
autorestart=true
priority=220
depends_on=xvfb
stdout_logfile=/dev/null
stderr_logfile=/dev/null
redirect_stderr=true

[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --vnc localhost:5901 --listen 6080 --heartbeat 30
autorestart=true
priority=230
depends_on=x11vnc
stdout_logfile=/dev/null
stderr_logfile=/dev/null
redirect_stderr=true

[program:nginx]
command=sudo nginx -g "daemon off;"
autorestart=true
priority=300
stdout_logfile=/dev/null
stderr_logfile=/dev/null
redirect_stderr=true

[program:code-server]
# Serves `/` on a local (non-Codespaces) sim. In Codespaces the wrapper
# exits 0 immediately — VS Code web is the Codespace itself and nginx maps
# `/` to 404 there (see entrypoint.sh root snippet selection).
command=/bin/bash -c "[ -n \"${CODESPACE_NAME:-}\" ] && exit 0; exec code-server --bind-addr 127.0.0.1:8080 --auth none --disable-telemetry --disable-update-check /home/physicar/physicar_ws"
user=physicar
environment=HOME="/home/physicar"
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=0
priority=310
stdout_logfile=/tmp/code-server.log
stderr_logfile=/tmp/code-server.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
redirect_stderr=true

[program:gz_websocket]
command=gz launch /opt/physicar/src/physicar-sim/websocket.gzlaunch
environment=GZ_CONFIG_PATH="/usr/share/gz",GZ_PARTITION="physicar",DISPLAY=":1"
autostart=false
autorestart=false
priority=400
stdout_logfile=/tmp/gz_websocket.log
stderr_logfile=/tmp/gz_websocket.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
redirect_stderr=true
startsecs=5

[program:sim_api]
command=python3 /opt/physicar/src/physicar-sim/sim_api.py
environment=DISPLAY=":1",GZ_CONFIG_PATH="/usr/share/gz",GZ_PARTITION="physicar",GALLIUM_DRIVER="llvmpipe",MESA_GL_VERSION_OVERRIDE="3.3"
autorestart=true
priority=410
startretries=999
stdout_logfile=/tmp/sim_api.log
stderr_logfile=/tmp/sim_api.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
redirect_stderr=true
startsecs=3

[program:physicar]
command=/bin/bash -c "/opt/physicar/src/physicar-ros/deploy/sim/physicar.sh"
autostart=true
autorestart=false
startsecs=0
priority=650
killasgroup=true
stopasgroup=true
stdout_logfile=/opt/physicar/userdata/physicar.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
stderr_logfile=/opt/physicar/userdata/physicar.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=0

[program:myapp]
command=/bin/bash -ic "[ -f /home/physicar/physicar_ws/myapp/run.sh ] || exit 0; exec bash /home/physicar/physicar_ws/myapp/run.sh"
directory=/opt/physicar
environment=PORT="5000",PYTHONUNBUFFERED="1"
# autostart=false so supervisord doesn't spawn this on a fresh sim with no
# run.sh — spawning would immediately create an empty run.log. entrypoint.sh
# starts it only when run.sh exists, and the web UI starts it on deploy.
# (Mirrors the device unit's ConditionPathExists=.../myapp/run.sh.)
autostart=false
autorestart=unexpected
startsecs=0
priority=700
stdout_logfile=/home/physicar/physicar_ws/myapp/run.log
stderr_logfile=/home/physicar/physicar_ws/myapp/run.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
stdout_logfile_backups=0
stderr_logfile_backups=0
redirect_stderr=true
stopsignal=TERM
killasgroup=true
stopasgroup=true

[program:memory_setup]
command=/bin/bash -c "[ -z \"${CODESPACE_NAME:-}\" ] && exit 0; exec bash /opt/physicar/src/physicar-ros/deploy/sim/memory-setup.sh"
autostart=true
autorestart=false
startsecs=0
priority=800
stdout_logfile=/tmp/memory-setup.log
stderr_logfile=/tmp/memory-setup.log
redirect_stderr=true

[program:app_browser]
command=/bin/bash -c "exec bash /opt/physicar/src/physicar-ros/deploy/sim/app-browser.sh"
autostart=true
autorestart=unexpected
exitcodes=0
priority=250
stdout_logfile=/tmp/app-browser.log
stderr_logfile=/tmp/app-browser.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
killasgroup=true
